# systemd/gpu-slack-runner.timer.template
[Unit]
Description=Run GPU Slack vLLM Runner every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
AccuracySec=5min
Persistent=true
Unit=gpu-slack-runner.service

[Install]
WantedBy=timers.target
